# systemd/quantum-pirate-miner.service
[Unit]
Description=Jones Quantum Pirate Miner - Headless 24/7 Operation with SphinxOSIIT Oracle
After=network.target network-online.target sphinxos.service
Wants=network-online.target
PartOf=sphinxos.service

[Service]
Type=simple
User=sphinxos
Group=sphinxos
WorkingDirectory=/opt/sphinxos/Sphinx_OS
Environment="PATH=/opt/sphinxos/Sphinx_OS/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/sphinxos/Sphinx_OS"
Environment="SDL_VIDEODRIVER=dummy"
Environment="SDL_AUDIODRIVER=dummy"
ExecStart=/opt/sphinxos/Sphinx_OS/venv/bin/python3 quantum_pirate_miner_headless.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-pirate-miner

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/sphinxos/Sphinx_OS /var/log
ProtectHome=true

# Resource limits
MemoryLimit=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
